# package/etc/init.d/start.sh
#! /bin/bash
# chkconfig: 2345 10 90
echo "ZMJ Service starting"
chmod +x /zmj/work/docker-hello
#/zmj/work/docker-hello
/zmj/work/XDtest

